# systemd-user/.config/systemd/user/battery_monitor_script.service
[Unit]
Description=Notify when battery is full or low

[Service]
ExecStart=/home/abubakr/.local/share/my_scripts/battery_popup.sh
