# systemd-user/.config/systemd/user/system_update_reminder.service
[Unit]
Description=Notify the user to update the system and Neovim

[Service]
ExecStart=/home/abubakr/.local/share/my_scripts/update_reminder.sh
